# crontab
# m h  dom mon dow   command
* * * * * /home/eliot/check_conf.sh
